# crontab
5 0 * * *   root cd /solar && /usr/bin/python3 /solar/multisma2/multisma2.py 2>&1
# Don't remove the empty line at the end of this file. It is required to run the cron job
